# backend/systemd/hantatracker.timer
[Unit]
Description=Run HantaTracker fetch every 5 minutes
Requires=hantatracker.service

[Timer]
# Cada 5 minutos (:00, :05, :10, :15, :20, :25, :30, :35, :40, :45, :50, :55)
OnCalendar=*-*-* *:00/5:00
# Si la Pi estaba apagada cuando tocaba, dispara al arrancar
Persistent=true
# Aleatoriza ±15s para no martillear servidores en el segundo exacto
RandomizedDelaySec=15
Unit=hantatracker.service

[Install]
WantedBy=timers.target
