# backend/systemd/hantatracker.timer
[Unit]
Description=Run HantaTracker fetch every 10 minutes
Requires=hantatracker.service

[Timer]
# Cada 10 minutos (:00, :10, :20, :30, :40, :50)
OnCalendar=*-*-* *:00/10:00
# Si la Pi estaba apagada cuando tocaba, dispara al arrancar
Persistent=true
# Aleatoriza ±30s para no martillear servidores en el segundo exacto
RandomizedDelaySec=30
Unit=hantatracker.service

[Install]
WantedBy=timers.target
